# curator4-1082331012.conf
[program:curator4-1082331012]
command=python3 /home/curator4/curator4-1082331012/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1082331012/test.err.log
stdout_logfile=/home/curator4/curator4-1082331012/test.out.log